# smtp-pro/ditmail.service
[Unit]
Description=DITMail SMTP Server (Haraka)
After=network.target redis.service
Wants=redis.service

[Service]
Type=simple
User=haraka
Group=haraka
WorkingDirectory=/home/ubuntu/DITMail/smtp-pro
ExecStart=sudo haraka -c /home/ubuntu/DITMail/smtp-pro
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ditmail

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/haraka /var/run/haraka /var/lib/haraka /var/mail

# Resource limits
LimitNOFILE=65536
LimitNPROC=32768

[Install]
WantedBy=multi-user.target
